# payjoin-mailroom/payjoin-mailroom.example.service
[Unit]
Description=Payjoin Mailroom
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/rust-payjoin/payjoin-mailroom
ExecStart=/root/.cargo/bin/payjoin-mailroom --config /root/rust-payjoin/payjoin-mailroom/config.toml
Restart=on-failure
RestartSec=5

# A public directory + OHTTP relay holds many concurrent file descriptors:
# inbound connections, long-poll waits, and OHTTP bootstrap tunnels (2 fds
# each). The systemd default soft limit (1024) is exhausted under load, making
# accept() fail with EMFILE ("Too many open files"). soft:hard form — 65536 is
# the enforced guardrail (well above peak); the 524288 hard ceiling stays
# available to raise live (prlimit) without a redeploy.
LimitNOFILE=65536:524288

[Install]
WantedBy=multi-user.target
